# supervisord_add.conf
[supervisord]
nodaemon=true

[program:apache]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true

[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
